# udev_keyboard/20-usb-keyboards.rules
ACTION=="add", ATTRS{idVendor}=="1006", ATTRS{idProduct}=="0022", RUN+="~/conf/udev_keyboard/udev.sh", OWNER="jerry"
